# system-packages/etc/supervisor/available.d/dnsmasq.conf
[program:dnsmasq]
command=/opt/redgate/system-packages/sbin/dnsmasq -d -C /opt/redgate/system-packages/etc/dnsmasq.conf
process_name=%(program_name)s

